# files/glinet/etc/hotplug.d/iface/30-adblock
#!/bin/sh

if [ "$ACTION" = "ifup" ];
	then if [ "$INTERFACE" = "wan" ]; then
		logger -t Adblock "Starting Adblock due to $ACTION of $INTERFACE ($DEVICE)"
		sh /etc/adblock.sh
	fi
fi